# etc/systemd/system/gammu.timer
[Unit]
Description=Run Gammu 40s later

[Timer]
OnBootSec=40
Unit=gammu-smsd.service

[Install]
WantedBy=multi-user.target
